# config/systemd/user/ipfs.service
# container-6ad50fe3dcc7f1c70275adca1106fa5f0c19c35af605a107ff4fd13041d4758e.service
# autogenerated by Podman 1.8.1-dev
# Wed Mar 25 04:21:01 -03 2020

[Unit]
Description=Podman container-6ad50fe3dcc7f1c70275adca1106fa5f0c19c35af605a107ff4fd13041d4758e.service
Documentation=man:podman-generate-systemd(1)

[Service]
Restart=on-failure
ExecStart=/usr/bin/podman start 6ad50fe3dcc7f1c70275adca1106fa5f0c19c35af605a107ff4fd13041d4758e
ExecStop=/usr/bin/podman stop -t 10 6ad50fe3dcc7f1c70275adca1106fa5f0c19c35af605a107ff4fd13041d4758e
PIDFile=/run/user/1000/containers/overlay-containers/6ad50fe3dcc7f1c70275adca1106fa5f0c19c35af605a107ff4fd13041d4758e/userdata/conmon.pid
KillMode=none
Type=forking

[Install]
WantedBy=default.target
