# config/systemd/user/tor.service
# container-4768b0215930a7a703d8e7723f98b5e1b110371dee0072bfa1df60b932cbf365.service
# autogenerated by Podman 1.8.1-dev
# Wed Mar 25 04:05:07 -03 2020

[Unit]
Description=Podman container-4768b0215930a7a703d8e7723f98b5e1b110371dee0072bfa1df60b932cbf365.service
Documentation=man:podman-generate-systemd(1)

[Service]
Restart=on-failure
ExecStart=/usr/bin/podman start 4768b0215930a7a703d8e7723f98b5e1b110371dee0072bfa1df60b932cbf365
ExecStop=/usr/bin/podman stop -t 10 4768b0215930a7a703d8e7723f98b5e1b110371dee0072bfa1df60b932cbf365
PIDFile=/run/user/1000/containers/overlay-containers/4768b0215930a7a703d8e7723f98b5e1b110371dee0072bfa1df60b932cbf365/userdata/conmon.pid
KillMode=none
Type=forking

[Install]
WantedBy=default.target
